# root/etc/s6-overlay/s6-rc.d/init-adguardhome-sync-config/up
/etc/s6-overlay/s6-rc.d/init-adguardhome-sync-config/run
